# docker/rootfs/usr/local/lib/swipl/doc/packages/examples/http/upstart-script.conf
# demo - SWI-Prolog demo server
#
# The SWI-Prolog demo server

description     "SWI-Prolog demo server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 5 60
umask 022

console log
chdir /home/swipl/src/demo

script
  export LANG=en_US.utf8
  /usr/bin/swipl daemon.pl --no-fork --port=80 --user=www-data --pidfile=/var/run/demo.pid --workers=16 --syslog=demo
end script
