# docker/rootfs/etc/supervisor/conf.d/supervisord_logicmoo.conf
# Hosts Logicmoo HTTP REPL at port 4123
[program:4123]
priority=10
user=root
directory=/opt/logicmoo_workspace/prologmud_server
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=/usr/local/bin/butterfly.server.py --port=%(program_name)s    --debug --log-to-stderr --i-hereby-declare-i-dont-want-any-security-whatsoever --unsecure --host=0.0.0.0 --cmd="bash --login -i /opt/logicmoo_workspace/bin/lmoo-bfly %(program_name)s"
#--cmd="screen -x -r LogicmooServer"

[program:flair_4095]
priority=10
user=root
directory=/opt/logicmoo_workspace/packs_sys/logicmoo_nlu/prolog/logicmoo_nlu
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=python3 ./parser_flair.py -port 4095


[program:spacy_4096]
priority=10
user=root
directory=/opt/logicmoo_workspace/packs_sys/logicmoo_nlu/prolog/logicmoo_nlu
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=python3 ./parser_spacy.py -port 4096

[program:allen_srl_4097]
priority=10
user=root
directory=/opt/logicmoo_workspace/packs_sys/logicmoo_nlu/prolog/logicmoo_nlu
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=python3 ./parser_allen_srl.py -port 4097

# 


# Hosts Logicmoo HTTP REPL at port 4125
[program:shrdlu]
priority=10
user=root
directory=/opt/logicmoo_workspace/packs_web/shrdlu
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=bash ./build
#--cmd="screen -x -r LogicmooServer"


# Hosts Logicmoo HTTP REPL at port 4125
[program:4125]
priority=10
user=root
directory=/opt/logicmoo_workspace/prologmud_server
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL
command=/usr/local/bin/butterfly.server.py --port=%(program_name)s    --debug --log-to-stderr --i-hereby-declare-i-dont-want-any-security-whatsoever --unsecure --host=0.0.0.0 --cmd="emacs -nw run_mud_server.pl"
#--cmd="screen -x -r LogicmooServer"


# Hosts Logicmoo MUD 4000-4001
[program:butterfly_relay]
process_name=relay_%(process_num)s
numprocs=5
command=/usr/local/bin/butterfly.server.py  --port=410%(process_num)s --debug --log-to-stderr --i-hereby-declare-i-dont-want-any-security-whatsoever --unsecure --host=0.0.0.0 --cmd="bash --login -i /opt/logicmoo_workspace/bin/lmoo-telnet 400%(process_num)s"
priority=10
user=root
directory=/opt/logicmoo_workspace/prologmud_server
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL

# Hosts /bin/login at port 4122
[program:22]
command=/usr/local/bin/butterfly.server.py --debug --i-hereby-declare-i-dont-want-any-security-whatsoever --unsecure --host=0.0.0.0 --port=41%(program_name)s --cmd="/bin/login"
priority=10
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL


[program:LogicmooServerLoop]
priority=10
user=root
directory=/opt/logicmoo_workspace
command=/opt/logicmoo_workspace/LogicmooServerLoop.sh
autostart = true
autorestart = true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL

# Hosts Logicmoo HTTPS REPL at port 14123
[program:apache2]
command=apache2ctl -e info  -c "ErrorLog /dev/stdout" -DFOREGROUND
priority=997
autostart=true
autorestart=true
startretries=1
startsecs=1
redirect_stderr=true
stdout_logfile=/var/log/myapache.out.log
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
user=root
killasgroup=true
stopasgroup=true

# Hosts Stanford HTTP at port 4090
[program:stanford]
priority=10
directory=/opt/logicmoo_workspace/prologmud_server/eclipse-workspace/stanford_corenlp
command=bash ./LogicmooStanford.sh
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL


[program:charniak]
priority=10
directory=/opt/logicmoo_workspace/prologmud_server/eclipse-workspace/bllip-parser
command=bash ./LogicmooCharniak.sh
# Now ran from SWI-Prolog
autostart=false
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL


# Hosts Eggdrop Telnet at port 3334
[program:eggdrop]
priority=10
user=prologmud_server
directory=/opt/logicmoo_workspace/packs_sys/eggdrop/conf
command=bash ./StartEggdrop.sh 
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
stopsignal=SIGKILL


#screen -mdS "jetty4077"
#screen -S jetty4077 -p0 -X stuff "$DIR0/prologmud_server/AceWiki/run-webapps-with-jettyrunner.sh\r"
# Hosts AceWiki 4077
[program:jetty4077]
directory=/opt/logicmoo_workspace/prologmud_server/eclipse-workspace/AceWiki/
command=bash ./run-webapps-with-jettyrunner.sh
priority=5
environment = PYTHONUNBUFFERED=1
autostart=false
autorestart=true
startsecs=4
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
killasgroup=true
stopasgroup=true
